# deploy/systemd/agent-radio-liquidsoap.service
[Unit]
Description=Agent Radio - Liquidsoap Audio Engine
After=icecast2.service
Wants=icecast2.service

[Service]
Type=simple
# TODO: Change 'radioagent' to the user that will run Agent Radio on your system
User=radioagent
WorkingDirectory=/opt/agent-radio
ExecStart=/usr/bin/liquidsoap /opt/agent-radio/radio.liq
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
